# zwave_me_server/rootfs/etc/services.d/z-way-server/finish
#!/usr/bin/with-contenv bashio

bashio::log.info "Stopping z-way-server..."

# Check if the PID file exists
if [ -f /var/run/z-way-server.pid ]; then
    # Fetch the PID from the PID file
    PID=$(cat /var/run/z-way-server.pid)

    # Kill the process using the PID
    kill $PID

    # Wait for the process to terminate
    wait $PID

    # Remove the PID file
    rm /var/run/z-way-server.pid

    bashio::log.info "z-way-server has been stopped successfully!"
else
    bashio::log.warning "z-way-server is not running or PID file is missing."
fi
